# crontab
*/10 * * * * ~/bin/rtcontrol is_ghost=yes --delete --yes
