# conf/crontab
*/10 * * * * ~/bin/rtcontrol is_complete=yes --delete --yes
